# 1.0/container-files/etc/services.d/sympa-task/run
#!/bin/execlineb -P
/usr/sbin/task_manager.pl
