# container-files/etc/services.d/sympa/run
#!/bin/execlineb -P

# s6-svwait -t 5000 -u /var/run/s6/services/syslog-ng
foreground { sleep 10 } foreground { /usr/sbin/presympa } /usr/sbin/sympa_msg.pl --foreground
