# 1.0/container-files/etc/services.d/sympa-bounce/run
#!/bin/execlineb -P
/usr/sbin/bounced.pl
