# 1.0/container-files/etc/services.d/_rsyslog/run
#!/bin/execlineb -P
/usr/sbin/rsyslogd -n
