# 1.0/container-files/etc/services.d/sympa/run
#!/bin/execlineb -P
/usr/sbin/presympa.sh && /usr/sbin/sympa_msg.pl --foreground
