# container-files/etc/services.d/sympa-bounce/run
#!/bin/execlineb -P

# s6-svwait -t 5000 -u /var/run/s6/services/sympa
foreground { sleep 15 } /usr/sbin/bounced.pl --foreground
