# container-files/etc/services.d/postfix/run
#!/bin/execlineb -P
# s6-svwait -t 5000 -U /var/run/s6/services/syslog-ng
foreground { /etc/syslog_wait } /usr/libexec/postfix/master -c /etc/postfix -d
