# container-files/etc/services.d/postfix/run
#!/bin/execlineb -P
/usr/libexec/postfix/master -c /etc/postfix -d
