# container-files/etc/services.d/sympa-bulk/run
#!/bin/execlineb -P
s6-svwait -t 5000 -U /var/run/s6/services/sympa
/usr/sbin/bulk.pl --foreground
