# 1.0/container-files/etc/services.d/sympa-archive/run
#!/bin/execlineb -P
/usr/sbin/archived.pl --foreground
