# 1.0/container-files/etc/services.d/sympa-msg/run
#!/bin/execlineb -P
foreground { /usr/sbin/presympa } /usr/sbin/sympa_msg.pl --foreground
